# deploy/services/comfyui.service
[Unit]
Description=ComfyUI AI Image Generation
Documentation=https://github.com/comfyanonymous/ComfyUI
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=nebula
Group=nebula
WorkingDirectory=/opt/nebula-command/ComfyUI

ExecStartPre=/bin/bash -c 'test -d /opt/nebula-command/ComfyUI/venv || (echo "Virtual environment not found" && exit 1)'

ExecStart=/opt/nebula-command/ComfyUI/venv/bin/python main.py \
    --listen 0.0.0.0 \
    --port 8188 \
    --preview-method auto \
    --use-pytorch-cross-attention

Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

Environment="PYTHONUNBUFFERED=1"
Environment="PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:512"
Environment="CUDA_VISIBLE_DEVICES=0"

StandardOutput=append:/opt/nebula-command/logs/comfyui.log
StandardError=append:/opt/nebula-command/logs/comfyui-error.log

LimitNOFILE=65536
LimitNPROC=4096

ProtectSystem=strict
ReadWritePaths=/opt/nebula-command/ComfyUI
ReadWritePaths=/opt/nebula-command/logs
ReadWritePaths=/opt/nebula-command/models
PrivateTmp=true
NoNewPrivileges=true

OOMScoreAdjust=100
Nice=10

MemoryMax=80%
TasksMax=256

[Install]
WantedBy=multi-user.target
