# deploy/services/nebula-discord-bot.service
[Unit]
Description=Nebula Command Discord Bot
Documentation=https://github.com/nebula-command/nebula-command
After=network-online.target postgresql.service nebula-dashboard.service
Wants=network-online.target

[Service]
Type=simple
User=nebula
Group=nebula
WorkingDirectory=/opt/nebula-command/services/discord-bot
ExecStartPre=/usr/bin/npm install --production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

Environment="NODE_ENV=production"
EnvironmentFile=-/opt/nebula-command/config/discord-bot.conf

LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nebula-discord-bot

WatchdogSec=60
NotifyAccess=main

[Install]
WantedBy=multi-user.target
